# Library/LaunchAgents/com.panda.watch-ice.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.panda.watch-ice</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/bo/.asdf/shims/bun</string>
    <string>/Users/bo/dev/tools/bun/scripts/watch-ice.js</string>
  </array>

  <key>KeepAlive</key>
  <true/>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>

